# home/Library/LaunchAgents/com.dbmrq.brewupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dbmrq.brewupdate.plist</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dbmrq/.brewupdate.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>604800</integer>

    <key>Nice</key>
    <integer>1</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/com.dbmrq.brewupdate.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/com.dbmrq.brewupdate.out</string>
</dict>
</plist>
